# rootfs/etc/services.d/amavisd/run
#!/usr/bin/execlineb -P

exec -c /usr/sbin/amavisd foreground
